# supervisorconf/supervisorord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB

[program:myapp]
command: /usr/bin/python -m ~/myapp/myapp
directory: /home/ubuntu
autostart: true
