# supervisorconf/supervisorord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB

[program:bookapp]
command: /usr/bin/python -m bookapp
directory: /home/ubuntu
autostart: true
